# root/etc/supervisord.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -display :0  -forever  -shared -o /var/log/x11vnc.log -rfbauth /home/app/.password2
autorestart=true